# memos/com.carlson.memo-transcribe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.carlson.memo-transcribe</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/carlson/.openclaw/workspace/memos/transcribe-cron.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/carlson/.openclaw/workspace/memos/transcribe-cron.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/carlson/.openclaw/workspace/memos/transcribe-cron.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
